# deploy/systemd/relay-node-updater.service
[Unit]
Description=Apply pending relay-node upgrade
Documentation=https://github.com/unix-relay/relay
# Don't run while relay-node is being deliberately stopped/started by us.
After=network-online.target
Wants=network-online.target
ConditionPathExists=/var/lib/relay-node/upgrade-request.json

[Service]
Type=oneshot
User=root
ExecStart=/usr/local/lib/relay-node/relay-node-updater
# Hard timeout: if updater hangs, give up so the path unit can re-trigger.
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
