# deploy/systemd/relay-node-updater.path
[Unit]
Description=Watch for relay-node upgrade requests
Documentation=https://github.com/unix-relay/relay

[Path]
PathExists=/var/lib/relay-node/upgrade-request.json
Unit=relay-node-updater.service

[Install]
WantedBy=multi-user.target
